# smartdial.conf
[program:smartdialAPI]
command=/root/smart/smartdial serve
directory=/root/smart
autostart=true
autorestart=true
stderr_logfile=/var/log/smartdial.err.log
stdout_logfile=/var/log/smartdial.out.log
